# rootfs/etc/inittab
# first run the system script file. See examples/inittab in the source of busybox for more information.

::sysinit:/etc/init.d/rcS

# example 1
#::askfirst:-/bin/sh
#::respawn:-/bin/sh

# example 2
# console form uboot bootargs. if you don't know the name of usart(such as eg.4.5.6), you can use this example.
console::askfirst:-/bin/sh

# example 3
#tty1::askfirst:-/bin/sh

# example 4
#s3c2410_serial2::askfirst:-/bin/sh

# example 5
#s3c2410_serial2::sysinit:/bin/login
#console::sysinit:/bin/login

# example 6
#s3c2410_serial2::respawn:/sbin/getty -L s3c2410_serial2 115200 vt100
#console::respawn:/sbin/getty -L console 115200 vt100

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

